# launchd/com.jarvis.ai-futures-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jarvis.ai-futures-digest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/terrybyrd/bin/jarvis-ai-futures-daily.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/terrybyrd/.local/share/jarvis/ai-futures/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/terrybyrd/.local/share/jarvis/ai-futures/logs/launchd-stderr.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/terrybyrd/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
